# sdcard/pi-kitchen/018-x11vnc-display/config/autostart/x11vnc.desktop
[Destop Entry]
Encoding=UTF-8
Type=Application
Name=X11VNC
Comment=
Exec=x11vnc -forever -usepw -display :0 -ultrafilexfer
StartupNotify=false
Terminal=false
Hidden=false
